# artifacts/persistence-methods/image_file_execution_options_hijack.yml
title: "Image File Execution Options Hijacking"
category: "persistence-methods"
description: "Process hijacking through Image File Execution Options debugger attachment and execution redirection"

paths:
  - "HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options"
  - "HKLM\\SOFTWARE\\WOW6432Node\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options"
  - "HKCU\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options"

details:
  what: |
    Image File Execution Options (IFEO) allows attaching debuggers to processes, modifying
    executable behavior, and controlling process execution parameters. Originally designed
    for debugging and compatibility, malware abuses this mechanism to hijack legitimate
    processes by setting malicious debuggers that run instead of the target application.

  forensic_value: |
    Advanced persistence technique that hijacks legitimate processes through debugger attachment.
    Shows unauthorized debugger attachments, process redirection, and sophisticated evasion
    techniques. Critical for detecting advanced malware persistence, accessibility tool hijacking
    (sticky keys attacks), and process replacement attacks. Can indicate privilege escalation
    attempts and system compromise.

  structure: |
    Executable names as subkeys containing Debugger value pointing to malicious executables,
    GlobalFlag for debugging options, VerifierDlls for application verifier, and various
    debugging parameters. When the hijacked process starts, the debugger runs instead with
    the original process as a parameter.

  examples:
    - "sethc.exe\\Debugger: C:\\Windows\\System32\\cmd.exe (Sticky Keys hijack)"
    - "taskmgr.exe\\Debugger: C:\\malware\\backdoor.exe"
    - "regedit.exe\\Debugger: powershell.exe -enc <base64_payload>"
    - "utilman.exe\\Debugger: C:\\temp\\shell.exe (Utility Manager hijack)"
    - "osk.exe\\Debugger: cmd.exe (On-Screen Keyboard hijack)"
    - "GlobalFlag: 0x200 (Enable heap validation)"
    - "VerifierDlls: malicious_verifier.dll"
    - "DisableHeapLookaside: 1"
    - "PageHeapFlags: 0x03"

  tools:
    - name: "Registry Explorer"
      url: "https://ericzimmerman.github.io/#!index.md"
      description: "Advanced registry browser for IFEO analysis and debugger detection"
    - name: "Autoruns"
      url: "https://docs.microsoft.com/en-us/sysinternals/downloads/autoruns"
      description: "Microsoft Sysinternals tool with IFEO enumeration and analysis"
    - name: "RegRipper"
      url: "https://github.com/keydet89/RegRipper3.0"
      description: "Registry analysis with Image File Execution Options plugins"
    - name: "Process Monitor"
      url: "https://docs.microsoft.com/en-us/sysinternals/downloads/procmon"
      description: "Real-time process monitoring for IFEO hijacking detection"
    - name: "IFEO Manager"
      description: "Tools for managing and analyzing Image File Execution Options"

limitations:
  - "IFEO debugger settings do NOT prove hijacked processes were actually executed"
  - "Debugger attachments don't indicate successful process redirection occurred"
  - "Registry entries don't show hijacked processes launched their debuggers"
  - "IFEO settings may be overridden by security software or policies"
  - "Process hijacking may fail due to missing debugger files or permission issues"
  - "Accessibility tool hijacking doesn't prove unauthorized access occurred"
  - "Configuration doesn't reveal frequency or success of hijacking attempts"

correlation:
  required_for_definitive_hijacking_proof:
    - "Process execution logs showing debugger execution instead of target processes"
    - "File system artifacts showing malicious debugger files at specified paths"
    - "Event logs showing process creation events for hijacked executables"
    - "Memory dumps showing debugger processes running in place of targets"
    
  strengthens_evidence:
    - "Registry modifications made by hijacked processes or debuggers"
    - "Network traffic initiated by malicious debuggers"
    - "File system changes indicating debugger activity"
    - "Application event logs showing failed target process launches"

metadata:
  windows_versions:
    - "Windows NT 3.5"
    - "Windows NT 4.0"
    - "Windows 2000"
    - "Windows XP"
    - "Windows Vista"
    - "Windows 7"
    - "Windows 8"
    - "Windows 8.1"
    - "Windows 10"
    - "Windows 11"
    - "Windows Server 2000"
    - "Windows Server 2003"
    - "Windows Server 2008"
    - "Windows Server 2012"
    - "Windows Server 2016"
    - "Windows Server 2019"
    - "Windows Server 2022"

  introduced: "Windows NT 3.5"
  criticality: "high"

  investigation_types:
    - "persistence-analysis"
    - "malware-analysis"
    - "privilege-escalation"
    - "incident-response"
    - "behavioral-analysis"

  tags:
    - "persistence"
    - "process-hijacking"
    - "debugger-abuse"
    - "advanced-malware"
    - "evasion"
    - "ifeo"
    - "sticky-keys"
    - "accessibility-hijack"
    - "process-replacement"

  references:
    - title: "Image File Execution Options"
      url: "https://docs.microsoft.com/en-us/previous-versions/windows/it-pro/windows-server-2003/cc786263(v=ws.10)"
      type: "official"
    - title: "Application Verifier"
      url: "https://docs.microsoft.com/en-us/windows-hardware/drivers/devtest/application-verifier"
      type: "official"
    - title: "MITRE ATT&CK: Image File Execution Options Injection"
      url: "https://attack.mitre.org/techniques/T1546/012/"
      type: "research"
    - title: "Sticky Keys Attack and IFEO"
      url: "https://www.sans.org/blog/malware-persistence-image-file-execution-options/"
      type: "research"

  retention:
    default_location: "SOFTWARE registry hive"
    persistence: "Survives reboots, highly persistent until explicitly removed"
    volatility: "IFEO settings persist across system updates and application reinstallation"

  related_artifacts:
    - "winlogon_userinit"
    - "com_objects"
    - "registry_run_keys"
    - "shell_extensions"
    - "app_compat_flags"

author:
  name: "Tonmoy Jitu"
  github: "tonmoy0010"
  x: "tonmoy0010"

contribution:
  date_added: "2025-01-15"
  last_updated: "2025-06-12"
  version: "3.0"
