# artifacts/persistence-methods/screensaver.yml
title: "Screensaver Persistence and Configuration"
category: "persistence-methods"
description: "Screensaver hijacking for persistence, execution triggers, and idle-time activation"

paths:
  - "HKCU\\Control Panel\\Desktop"
  - "HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\IniFileMapping\\Control.ini"
  - "HKCU\\Control Panel\\Screen Saver.*"

details:
  what: |
    Windows screensaver configuration including screensaver executable path, timeout settings,
    password protection, and screensaver-related security settings. Screensavers run with user
    privileges when activated and provide an execution mechanism triggered by system idle time.
    Configuration includes screensaver selection, display properties, and activation policies.

  forensic_value: |
    Malware can replace legitimate screensavers with malicious executables for persistence that
    triggers during idle periods. Shows unauthorized screensaver modifications that provide
    execution opportunities when systems are unattended. Can indicate sophisticated persistence
    mechanisms that activate during low-activity periods to avoid detection.

  structure: |
    Desktop settings include SCRNSAVE.EXE (screensaver path), ScreenSaveTimeOut (activation time
    in seconds), ScreenSaverIsSecure (password protection), ScreenSaveActive (enabled/disabled),
    and various screensaver-specific configuration parameters stored as REG_SZ values in
    Control Panel\\Desktop.

  examples:
    - "SCRNSAVE.EXE: C:\\Windows\\System32\\Mystify.scr"
    - "ScreenSaveTimeOut: 900 (15 minutes)"
    - "ScreenSaverIsSecure: 1 (password required to unlock)"
    - "ScreenSaveActive: 1 (screensaver enabled)"
    - "Malicious: SCRNSAVE.EXE: C:\\malware\\backdoor.scr"
    - "Persistence: SCRNSAVE.EXE: powershell.exe -WindowStyle Hidden -File C:\\temp\\payload.ps1"
    - "Hijacked: SCRNSAVE.EXE: C:\\Users\\user\\AppData\\Local\\evil.exe"
    - "ScreenSaveUsePassword: 0 (no password protection)"

  tools:
    - name: "Display Properties"
      description: "Windows Display Properties control panel for screensaver configuration"
    - name: "Registry Explorer"
      url: "https://ericzimmerman.github.io/#!index.md"
      description: "Advanced registry browser for screensaver configuration analysis"
    - name: "Autoruns"
      url: "https://docs.microsoft.com/en-us/sysinternals/downloads/autoruns"
      description: "Microsoft Sysinternals tool with screensaver enumeration"
    - name: "Screen Saver Settings"
      description: "Windows screensaver configuration and management utilities"
    - name: "ScreenSaverView"
      description: "Tools for analyzing and monitoring screensaver configurations"

limitations:
  - "Screensaver configuration does NOT prove screensaver was activated or executed"
  - "Executable path settings don't indicate successful screensaver launch"
  - "Timeout settings don't show system idle time was reached"
  - "Malicious screensaver paths don't prove unauthorized execution occurred"
  - "Screensaver activation may be prevented by user activity or power settings"
  - "Registry settings don't reveal frequency or duration of screensaver execution"
  - "Password protection settings don't indicate screensaver lock/unlock events"

correlation:
  required_for_definitive_execution_proof:
    - "Process execution logs showing screensaver program startup"
    - "File system artifacts showing screensaver executable files at configured paths"
    - "Windows Event Logs showing process creation for screensaver programs"
    - "User activity logs showing system idle periods matching screensaver timeouts"
    
  strengthens_evidence:
    - "Registry modifications made by executed screensaver programs"
    - "Network traffic initiated by malicious screensaver executables"
    - "File system changes indicating screensaver program activity"
    - "Power management logs showing display state changes"

metadata:
  windows_versions:
    - "Windows 95"
    - "Windows 98"
    - "Windows NT 4.0"
    - "Windows 2000"
    - "Windows XP"
    - "Windows Vista"
    - "Windows 7"
    - "Windows 8"
    - "Windows 8.1"
    - "Windows 10"
    - "Windows 11"

  introduced: "Windows 95"
  criticality: "low"

  investigation_types:
    - "persistence-analysis"
    - "malware-analysis"
    - "behavioral-analysis"
    - "incident-response"

  tags:
    - "persistence"
    - "screensaver"
    - "execution"
    - "desktop-hijacking"
    - "idle-execution"
    - "user-privileges"
    - "unattended-execution"
    - "display-properties"

  references:
    - title: "Screensaver Configuration"
      url: "https://docs.microsoft.com/en-us/previous-versions/windows/desktop/legacy/dd162397(v=vs.85)"
      type: "official"
    - title: "Desktop Window Manager"
      url: "https://docs.microsoft.com/en-us/windows/win32/dwm/dwm-overview"
      type: "official"
    - title: "Screensaver Malware Persistence"
      url: "https://www.sans.org/blog/screensaver-persistence/"
      type: "research"
    - title: "Alternative Persistence Mechanisms"
      url: "https://attack.mitre.org/techniques/T1547/"
      type: "research"

  retention:
    default_location: "NTUSER.DAT registry hive"
    persistence: "Survives reboots, persists per user profile until configuration changes"
    volatility: "Screensaver settings preserved across system updates and power cycles"

  related_artifacts:
    - "registry_run_keys"
    - "shell_extensions"
    - "scheduled_tasks"
    - "winlogon_userinit"
    - "user_profiles"

author:
  name: "Tonmoy Jitu"
  github: "tonmoy0010"
  x: "tonmoy0010"

contribution:
  date_added: "2025-01-15"
  last_updated: "2025-06-12"
  version: "3.0"
